# src/control_plane/bacon-control.service
[Unit]
Description=BACON-AI Control Plane
After=network.target mosquitto.service

[Service]
Type=simple
User=bacon
WorkingDirectory=/home/bacon/bacon-control-plane
Environment=BACON_DB_PATH=/home/bacon/bacon-control-plane/bacon.db
Environment=MQTT_BROKER=localhost
Environment=MQTT_PORT=1883
Environment=MEM0_API_KEY=m0-TqhAnGlqmJt3Gr6SaJ8z5iaPWYVcGxpiI67XC4Xo
ExecStart=/home/bacon/bacon-control-plane/venv/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
