# src/control_plane/bacon-feature-physics.service
[Unit]
Description=BACON-AI Control Plane (feature-physics)
After=network.target mosquitto.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bacon-ai/feature-physics
ExecStart=/opt/bacon-ai/feature-physics/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8006
Restart=always
RestartSec=5
Environment="BACON_ENV=feature-physics"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
